# polybot.service
[Unit]
Description=Polybot Telegram Service
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/PolybotService
ExecStart=/home/ec2-user/PolybotService/venv/bin/python bot.py
Restart=always

[Install]
WantedBy=multi-user.target